# core-master/core-master/core/roles/httpd/templates/etc/cron.d/apache_dos_protect
* * * * * root for i in {0..5}; do /srv/southbridge/bin/apache_dos_protect.sh >/dev/null 2>&1 && sleep 10; done;
