# core-master/core-master/core/roles/nginx/templates/etc/systemd/system/nginx.service.d/override.conf
[Service]
LimitNOFILE=100000
PrivateTmp=no
