# core-master/core-master/core/roles/postfix/templates/etc/cron.d/resend_mail
*/10 * * * * root /etc/postfix/resend_mail > /dev/null 2>&1
